# centos7-django/supervisord_djangoapp.ini
[program:djangoapp]
command=gunicorn djangoapp.wsgi:application --bind 127.0.0.1:8000
directory=/djangoapp
user=root
autostart=true
autorestart=true
stdout_logfile=/djangoapp/supervisord.log
redirect_stderr=true
